# lcdtempr
### BEGIN INIT INFO
# Provides: LCD - time / CPU temp / room temp
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: time / CPU temp / room temp
# Change path to your temprmonitor.py path
# Make script exetutable: sudo chmod +x /etc/init.d/lcdtempr
# Make the lcdtempr init script known to the system by using the 
# update-rc.d command: sudo update-rc.d lcdtempr defaults
### END INIT INFO
 
 
#! /bin/sh
# /etc/init.d/lcdtempr
 
 
export HOME
case "$1" in
    start)
        echo "Starting LCD"
        
<<<<<<< HEAD
    sudo python /your/path/to/rpimonitor_stats/temprmonitor.py &
=======
sudo python /your/path/to/rpimonitor_stats/temprmonitor.py &
>>>>>>> 66e2048b862be86c3b8ff749599273dab54fa4f3
    ;;
    stop)
        echo "Stopping LCD"
    LCD_PID=`ps auxwww | grep temprmonitor.py | 
head -1 | awk '{print $2}'`
    kill -9 $LCD_PID
    ;;
    *)
        echo "Usage: /etc/init.d/lcdtempr {start|stop}"
        exit 1
    ;;
esac
exit 0
